# infrastructure/ami/core.service
[Unit]
Description = some random service

[Service]
Type=simple
ExecStart=/usr/bin/srv

[Install]
WantedBy=multi-user.target  